# mc-spawn-agent.service
# Reference systemd unit (install.sh generates the real one per backend).
# Kept here for inspection; you don't install this file directly.
#
# install.sh writes either a system unit (root), a --user unit (rootless), an
# OpenRC script (Alpine), or a nohup launcher — all of them exec $DIR/run.sh,
# a 0600 launcher that exports CONTROL_URL/TOKEN/AGENT_STATE (so secrets stay
# out of the world-readable unit file and out of `ps`).
[Unit]
Description=mc-spawn agent
After=network-online.target docker.service
Wants=network-online.target

[Service]
# run.sh exports CONTROL_URL / TOKEN / AGENT_STATE, then execs python3 agent.py.
ExecStart=/opt/mc-spawn-agent/run.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
